# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[program:api]
command=node /app/api/dist/server.js
directory=/app/api
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=NODE_ENV="production",PORT="%(ENV_API_PORT)s",DATABASE_PATH="%(ENV_DATABASE_PATH)s",INTERNAL_TOKEN="%(ENV_INTERNAL_TOKEN)s",TRUSTED_ORIGINS="%(ENV_TRUSTED_ORIGINS)s",DEFAULT_UPLOAD_PATH="%(ENV_DEFAULT_UPLOAD_PATH)s",HOST_ROOT_PREFIX="%(ENV_HOST_ROOT_PREFIX)s"

[program:worker]
command=node /app/worker/dist/main.js
directory=/app/worker
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=NODE_ENV="production",INTERNAL_TOKEN="%(ENV_INTERNAL_TOKEN)s",API_URL="http://127.0.0.1:%(ENV_API_PORT)s",HOST_ROOT_PREFIX="%(ENV_HOST_ROOT_PREFIX)s"
startsecs=5

[program:web]
command=node /app/web/dist/opengallery/server/server.mjs
directory=/app/web
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=NODE_ENV="production",PORT="4200",PUBLIC_URL="http://127.0.0.1:%(ENV_API_PORT)s"
startsecs=5

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
startsecs=2
